# isnsd.socket
[Socket]
ListenStream=/var/run/isnsctl
PassCredentials=true

[Socket]
ListenStream=3205

[Install]
WantedBy=sockets.target
